# control/slackbot.service
[Unit]
Description=Harmony Slack Transcription Bot
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=aaron
WorkingDirectory=%h/Documents/harmony
EnvironmentFile=%h/Documents/harmony/.env
ExecStart=%h/Documents/harmony/.venv/bin/python control/slackbot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
