# control/slackbot.service
[Unit]
Description=Harmony Slack Transcription Bot
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
# Before installing: replace /home/aaron with the actual user's home directory.
User=aaron
WorkingDirectory=/home/aaron/Documents/harmony
EnvironmentFile=/home/aaron/Documents/harmony/.env
ExecStart=/home/aaron/Documents/harmony/.venv/bin/python control/slackbot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
