# compute/transcribe_api.service
[Unit]
Description=Harmony Transcription API (faster-whisper)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
# Before installing: replace /home/aaron with the actual user's home directory.
User=aaron
WorkingDirectory=/home/aaron/Documents/harmony
EnvironmentFile=/home/aaron/Documents/harmony/.env
ExecStart=/home/aaron/Documents/harmony/.venv/bin/uvicorn compute.transcribe_api:app \
    --host ${FORGE_TAILSCALE_IP} \
    --port 8765
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
